# deploy/synk-bot.service
[Unit]
Description=Synk Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/synk
EnvironmentFile=/root/synk/.env
ExecStart=/root/synk/.venv/bin/python main.py
Restart=always
RestartSec=15
StandardOutput=append:/root/synk/logs/bot.log
StandardError=append:/root/synk/logs/bot.log

[Install]
WantedBy=multi-user.target
